# etc/systemd/system/static-arp.service
[Unit]
Description="Static ARP configuration script"

[Service]
Type=oneshot
ExecStart=/etc/scripts/static-arp

[Install]
WantedBy=network.target
